# proj/pyapp/docker_config/supervisor_app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /project/docker_config/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
